# samples/dummy.service
[Service]
Type=notify
ExecStart=/usr/local/bin/dummy
Restart=no
WatchdogSec=3s
